# mediamtx/rootfs/etc/services.d/mediamtx/run
#!/command/with-contenv bashio
set -euo pipefail

CFG="/data/mediamtx.yml"
bashio::log.info "Starting MediaMTX with ${CFG}"
exec /usr/bin/mediamtx "${CFG}"
